# scripts/com.historyforge.daily-video.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.historyforge.daily-video</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/adamschmidt/The-History-Forge/scripts/run_daily_video_job.sh</string>
    </array>

    <!--
      Run every day at 07:00 local time.
      If the Mac is asleep at 07:00, launchd will fire the job
      as soon as the Mac wakes (thanks to StartCalendarInterval).
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/adamschmidt/The-History-Forge</string>

    <key>StandardOutPath</key>
    <string>/Users/adamschmidt/The-History-Forge/data/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/adamschmidt/The-History-Forge/data/launchd_stderr.log</string>

    <!-- Run missed job after wake — this is the key advantage over cron -->
    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
